# bin/init.d/shinken-broker
#!/bin/bash
#
# Shinken        Starts Shinken
#
#
# chkconfig: 345 98 02
# description: Start Shinken broker daemon

# Source function library.
if [ -f /etc/init.d/functions ]; then
        . /etc/init.d/functions
fi

SHINKEN_USER=shinken
SHINKEN_GROUP=shinken

SHINKEN_BIN=/usr/local/shinken/bin
SHINKEN_VAR=/usr/local/shinken/src/var/

start() {
        echo "Starting Shinken Poller daemon  ..."
        su - $SHINKEN_USER -c "$SHINKEN_BIN/launch_broker.sh"
#        chown -R $SHINKEN_GROUP:$NAGIOS_GROUP $SHINKEN_VAR
        echo_success 
        echo

}

stop() {
        echo "Stopping Shinken Poller daemon  ..."
        su - $SHINKEN_USER -c "$SHINKEN_BIN/stop_broker.sh"
        #sleep 3
        #count=$(ps aux | grep python2.6 | grep -v grep | wc -l)
        #if [ $count -gt 0 ]; then
        #    killall python2.6
        #fi
        echo_success 
        echo
 
}

case "$1" in
  restart)
        stop
        start
        exit 0
        ;;
  start)
        start
        exit 0
        ;;
  stop)
       stop
       exit 0
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart}"
        exit 2
esac

exit $?

