# bin/init.d/shinken-arbiter
#! /bin/sh

### BEGIN INIT INFO
# Provides:          shinken-arbiter
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Shinken arbiter daemon
# Description:       Shinken is a monitoring tool and the Arbiter
#                    is one of its daemon. This one reads the configuration,
#                    cuts it into parts and dispatches it. Then it waits
#                    for orders from the users to dispatch them too.
### END INIT INFO

# Author: Gabes Jean <naparuba@gmail.com>
#         Olivier LI-KIANG-CHEONG <lkco@gezen.fr>


NAME="shinken-arbiter"

curdir=$(dirname "$0")

export SHINKEN_MODULE_FILE="$NAME"  ## for 'shinken' init script to see that it's called by us

"$curdir/shinken" $@ arbiter

exit $?

