# bin/init.d/shinken-reactionner
#! /bin/sh

### BEGIN INIT INFO
# Provides:          shinken-reactionner
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Shinken reactionner daemon
# Description:       Shinken is a monitoring tool and the Reactionner
#                    is one of its daemon. This one gets the configuration from the arbiter
#                    His purpose is to get the actually do the actions like sending an email
#                    ordered by the schedulers specified in the configuration
### END INIT INFO

### BEGIN INIT INFO Redhat
# Shinken        Starts Shinken Reactionner
#
# chkconfig: 345 99 01
# description: Start Shinken reactionner daemon
### END INIT INFO

# Author: Gabes Jean <naparuba@gmail.com>
#         Olivier LI-KIANG-CHEONG <lkco@gezen.fr>


NAME="shinken-reactionner"

curdir=$(dirname "$0")

export SHINKEN_MODULE_FILE="$NAME"  ## for 'shinken' init script to see that it's called by us

"$curdir/shinken" $@ reactionner

exit $?
