# debian/shinken-common.postrm
#!/bin/bash
#DEBHELPER#

set -e

update-rc.d -f shinken remove
