# debian/shinken-common.postrm
#!/bin/bash
#DEBHELPER#

set -e

update-rc.d -f shinken remove

update-rc.d -f shinken-arbiter remove
update-rc.d -f shinken-scheduler remove
update-rc.d -f shinken-poller remove
update-rc.d -f shinken-broker remove
update-rc.d -f shinken-receiver remove
update-rc.d -f shinken-reactionner remove

