# debian/shinken-poller.postrm
#!/bin/bash

set -e

update-rc.d -f shinken-poller remove
