# debian/shinken-poller.postinst
#!/bin/bash
#DEBHELPER#

set -e

mkdir -p /var/log/shinken/
chown -R shinken:shinken /var/log/shinken/

update-rc.d shinken-poller defaults
