# bin/rc.d/shinken-arbiter
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: shinken_arbiter
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="shinken_arbiter"
rcvar="shinken_arbiter_enable"

shinken_arbiter_configfile="/usr/local/etc/shinken/shinken.cfg"
command="/usr/local/bin/shinken-arbiter"
command_interpreter="/usr/local/bin/python2.7"
command_args="-d -c ${shinken_arbiter_configfile} > /dev/null 2>&1"
pidfile="/var/run/shinken/arbiterd.pid"
shinken_error_log="/tmp/shinken-arbiter.log"

restart_precmd="shinken_checkconfig"
configtest_cmd="shinken_checkconfig"

required_files="${shinken_arbiter_configfile}"
extra_commands="configtest"

load_rc_config "${name}"

[ -z "${shinken_arbiter_enable}" ] && shinken_arbiter_enable="NO"

shinken_checkconfig() {
	echo -n "Performing sanity check on shinken configuration: "
        ${command} -v -c ${shinken_arbiter_configfile} >${shinken_error_log} 2>&1
        if [ $? -ne 0 ]; then
		echo "FAILED, check ${shinken_error_log}"
		return 1
	else
		echo "OK"
		rm ${shinken_error_log}
	fi
}

run_rc_command "$1"
