# bin/systemd/shinken-arbiter.service
[Unit]
Description="Shinken - a modern, Nagios compatible monitoring framework"
Documentation=http://www.shinken.io/
Requires=network-online.target
After=network-online.target

[Service]
Type=forking
User=shinken
Group=shinken
EnvironmentFile=-/etc/default/shinken-arbiter
ExecStart=shinken-arbiter -d -c $SHINKENCFG -c $SHINKENSPECIFICCFG $DEBUGCMD
KillMode=mixed

[Install]
WantedBy=multi-user.target
