# bin/systemd/shinken-receiver.service
[Unit]
Description="Shinken - a modern, Nagios compatible monitoring framework"
Documentation=http://www.shinken.io/
Requires=network-online.target
After=network-online.target

[Service]
Type=forking
User=shinken
Group=shinken
EnvironmentFile=-/etc/default/shinken-receiver
ExecStart=shinken-receiver -d -c $RECEIVERCFG $DAEMON_ARGS
KillMode=mixed

[Install]
WantedBy=multi-user.target
